# src/securitas-wall-host/tmp/src/etc/network/if-up.d/iptables
#!/bin/sh
iptables-restore < /etc/iptables.ipv4
ip6tables-restore < /etc/iptables.ipv6