# src/securitas-wall-host/tmp/src/etc/rc.local
systemctl daemon-reload
/etc/init.d/clamav-daemon  restart

exit 0
